# cmdeploy/src/cmdeploy/filtermail/filtermail-incoming.service.j2
[Unit]
Description=Incoming Chatmail Postfix before queue filter 

[Service]
ExecStart={{ bin_path }} {{ config_path }} incoming
Restart=always
RestartSec=30
User=vmail

[Install]
WantedBy=multi-user.target
