# cmdeploy/src/cmdeploy/acmetool/acmetool-reconcile.service
[Unit]
Description=Renew TLS certificates with acmetool
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/acmetool --batch reconcile

